# ops/systemd/absurd-cleanup.timer
[Unit]
Description=Daily ABSURD Recovery Cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
