# ops/systemd/contimg-stream.service
[Unit]
Description=DSA-110 Continuum Streaming Converter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/data/dsa110-contimg/ops/systemd/contimg.env
WorkingDirectory=/data/dsa110-contimg
ExecStart=/bin/bash /data/dsa110-contimg/scripts/ops/casa/casa_wrapper.sh /opt/miniforge/envs/casa6/bin/python -m dsa110_contimg.conversion.streaming.streaming_converter \
  --input-dir ${CONTIMG_INPUT_DIR} \
  --output-dir ${CONTIMG_OUTPUT_DIR} \
  --queue-db ${CONTIMG_QUEUE_DB} \
  --registry-db ${CONTIMG_REGISTRY_DB} \
  --scratch-dir ${CONTIMG_SCRATCH_DIR} \
  --log-level ${CONTIMG_LOG_LEVEL} \
  --use-subprocess \
  --expected-subbands ${CONTIMG_EXPECTED_SUBBANDS} \
  --chunk-duration ${CONTIMG_CHUNK_MINUTES} \
  --monitoring \
  --monitor-interval ${CONTIMG_MONITOR_INTERVAL}
Restart=always
RestartSec=5
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=4
NoNewPrivileges=true
LimitNOFILE=1048576
StandardOutput=append:/data/dsa110-contimg/state/logs/contimg-stream.out
StandardError=append:/data/dsa110-contimg/state/logs/contimg-stream.err

[Install]
WantedBy=multi-user.target
