# ops/systemd/contimg-pointing-monitor.service
[Unit]
Description=DSA-110 Pointing Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
EnvironmentFile=/data/dsa110-contimg/ops/systemd/contimg.env
Environment=PYTHONPATH=/data/dsa110-contimg/src/dsa110_contimg/src
Environment=CASACORE_NOROOT=true
WorkingDirectory=/data/dsa110-contimg
ExecStart=/bin/bash /data/dsa110-contimg/scripts/ops/casa/casa_wrapper.sh /opt/miniforge/envs/casa6/bin/python -m dsa110_contimg.pointing.monitor \
    ${CONTIMG_INPUT_DIR} \
    ${PIPELINE_PRODUCTS_DB} \
    --status-file ${PIPELINE_RUN_DIR}/pointing-monitor-status.json \
    --health-check-interval 300 \
    --log-level ${CONTIMG_LOG_LEVEL}
Restart=always
RestartSec=60
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
