# ops/systemd/contimg-absurd-worker@.service
[Unit]
Description=DSA-110 ABSURD Workflow Worker %i
Documentation=file:///data/dsa110-contimg/docs/guides/absurd-deployment-plan.md
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/data/dsa110-contimg/backend

# Environment
EnvironmentFile=/data/dsa110-contimg/ops/systemd/contimg.env
Environment=ABSURD_WORKER_ID=worker-%i

# Activate conda environment and run worker
ExecStart=/bin/bash -c 'source /opt/miniforge/etc/profile.d/conda.sh && \
    conda activate casa6 && \
    python -m dsa110_contimg.absurd.worker'

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Resource limits
MemoryMax=64G
CPUQuota=800%

# Logging
StandardOutput=append:/data/dsa110-contimg/state/logs/absurd-worker-%i.log
StandardError=append:/data/dsa110-contimg/state/logs/absurd-worker-%i.log

# Watchdog for health monitoring
WatchdogSec=120

[Install]
WantedBy=multi-user.target
