# ops/systemd/contimg-docs.service
[Unit]
Description=DSA-110 Continuum Imaging Documentation Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
EnvironmentFile=/data/dsa110-contimg/ops/systemd/contimg.env
WorkingDirectory=/data/dsa110-contimg
ExecStart=/opt/miniforge/envs/casa6/bin/python -m mkdocs serve --dev-addr=127.0.0.1:${CONTIMG_DOCS_PORT}
Restart=always
RestartSec=5
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
