# ops/systemd/dsa110-disk-monitor.service
[Unit]
Description=DSA-110 Disk Usage Monitor
Documentation=file:///data/dsa110-contimg/internal/docs/dev-notes/analysis/INTEGRATION_GUIDE.md
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=dsa
Group=dsa
WorkingDirectory=/data/dsa110-contimg/src
Environment="PYTHONPATH=/data/dsa110-contimg/src"
Environment="DSA110_MONITOR_INTERVAL_DISK=600"
Environment="DSA110_ALERT_DISK_USAGE_PCT=80"
ExecStart=/opt/miniforge/envs/casa6/bin/python -m dsa110_contimg.qa.disk_usage_monitoring_daemon
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dsa110-disk-monitor

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ReadWritePaths=/data/dsa110-contimg/qa_outputs /tmp/dsa110_monitoring_plots
ProtectHome=false

[Install]
WantedBy=multi-user.target
