# scripts/services/casa-log-daemon-inotify.service
[Unit]
Description=CASA Log Daemon (inotifywait) - Monitors and moves casa-*.log files recursively
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/data/dsa110-contimg
ExecStart=/bin/bash /data/dsa110-contimg/scripts/casa_log_daemon_inotify.sh
Restart=always
RestartSec=10
# Restart on failure, but wait longer if it keeps failing
StartLimitInterval=300
StartLimitBurst=5
# Note: WatchdogSec removed - bash scripts don't send watchdog heartbeats
# The separate casa-log-daemon-watchdog.timer handles health monitoring instead
# Log to file for easier debugging
StandardOutput=append:/data/dsa110-contimg/state/logs/casa-log-daemon.out
StandardError=append:/data/dsa110-contimg/state/logs/casa-log-daemon.err

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/data/dsa110-contimg

[Install]
WantedBy=multi-user.target

