# systemd/dsa110-redis.service
# Redis Service for DSA-110 Pipeline
# Systemd service configuration

[Unit]
Description=Redis for DSA-110 Pipeline
Documentation=https://redis.io/
After=network.target
Wants=network.target

[Service]
Type=notify
User=redis
Group=redis
WorkingDirectory=/var/lib/redis
Environment=REDIS_CONF=/etc/redis/redis.conf
ExecStart=/usr/bin/redis-server /etc/redis/redis.conf --supervised systemd
ExecStop=/usr/bin/redis-cli shutdown
ExecReload=/bin/kill -USR2 $MAINPID
TimeoutStopSec=0
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/redis /var/log/redis
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768
MemoryMax=2G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dsa110-redis

[Install]
WantedBy=multi-user.target
