# pipeline/scripts/dsa110-streaming-converter.service
[Unit]
Description=DSA-110 Streaming Converter
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/data/dsa110-contimg
Environment=PATH=/opt/conda/envs/casa6/bin:/usr/local/bin:/usr/bin:/bin
Environment=CONDA_DEFAULT_ENV=casa6
Environment=CONDA_PREFIX=/opt/conda/envs/casa6

# Main command - customize paths and options as needed
ExecStart=/opt/conda/envs/casa6/bin/python streaming/streaming_converter.py \
    --input-dir /data/incoming_data \
    --output-dir /data/output/ms \
    --scratch-dir /data/scratch \
    --checkpoint-dir /data/checkpoints \
    --chunk-duration 5.0 \
    --omp-threads 4 \
    --monitor-interval 60 \
    --log-level INFO \
    --cleanup-temp
# Add "--no-monitoring" to ExecStart above if you prefer to disable monitoring.
# Adjust "--chunk-duration" to change the data window size (default: 5.0 minutes).

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits (adjust based on system)
LimitNOFILE=65536
LimitNPROC=32768

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dsa110-streaming-converter

[Install]
WantedBy=multi-user.target
